# deploy/scrabble.rc
#!/bin/sh

# PROVIDE: scrabble
# REQUIRE: LOGIN postgresql
# KEYWORD: shutdown

. /etc/rc.subr

name="scrabble"
rcvar="${name}_enable"

load_rc_config $name

: ${scrabble_enable:="NO"}
: ${scrabble_dir:="/opt/scrabble"}
: ${scrabble_env_file:="/opt/scrabble/.env"}

pidfile="${scrabble_dir}/${name}.pid"
procname="/usr/local/bin/node"
command="/usr/sbin/daemon"
command_args="-u scrabble -p ${pidfile} -o /var/log/${name}.log ${procname} ${scrabble_dir}/packages/server/dist/index.js"

start_precmd="${name}_prestart"

scrabble_prestart()
{
    if [ -f "${scrabble_env_file}" ]; then
        set -a
        . "${scrabble_env_file}"
        set +a
    fi
    cd "${scrabble_dir}"
}

run_rc_command "$1"
